# .docker/supervisor/app.ini
[program:nginx]
command=nginx
autostart=true
autorestart=true
priority=15
stdout_events_enabled=true
stderr_events_enabled=true

[program:php-fpm]
command=php-fpm7 -F
autostart=true
autorestart=true
priority=15
stdout_events_enabled=true
stderr_events_enabled=true

[program:mariadb]
command=mysqld --user=me
autostart=true
autorestart=true
priority=15
stdout_events_enabled=true
stderr_events_enabled=true